# utils/xtreamui.cron
# Xtream-UI
# A cada 3 horas, limpa os logs (rodando como root).
0 */3 * * * root sh -c ': > /home/xtreamcodes/iptv_xtream_codes/logs/main.error.log && : > /home/xtreamcodes/iptv_xtream_codes/logs/main.access.log'
