# utils/crons/xtreamui.cron
# Xtream-UI
0 */3 * * * root /bin/mkdir -p /home/xtreamcodes/iptv_xtream_codes/logs && /usr/bin/truncate -s 0 /home/xtreamcodes/iptv_xtream_codes/logs/main.error.log /home/xtreamcodes/iptv_xtream_codes/logs/main.access.log

# Xtream-UI
#* * * * * root /home/xtreamcodes/iptv_xtream_codes/php/bin/php -q /home/xtreamcodes/iptv_xtream_codes/crons/pid_monitor.php >/dev/null 2>&1